# apps/pod/deploy/podex-pod.service
[Unit]
Description=Podex Pod kiosk
After=network-online.target
Wants=network-online.target
Conflicts=getty@tty1.service

[Service]
Type=simple
User=podex
Group=podex
SupplementaryGroups=input render video gpio
WorkingDirectory=/opt/podex-pod
EnvironmentFile=/etc/podex-pod.env
Environment=PODEX_SIMULATOR=0
Environment=XDG_RUNTIME_DIR=/run/podex-pod
RuntimeDirectory=podex-pod
RuntimeDirectoryMode=0700
PAMName=login
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
StandardInput=tty
StandardOutput=journal
ExecStart=/usr/bin/cage -s -- /opt/podex-pod/.venv/bin/python -m podex_pod
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
